# kindle/etc/upstart/lazybones.conf
# upstart config file for kindle-lazybones

start on started lab126
stop on stopping system

exec /usr/bin/python /opt/lazybones/srv.py >> /var/log/lazybones.log 2>&1
